# mindfulity.cron
# /etc/cron.d/mindfulity

# run every quarter-hour
01,16,31,46 *     * * *     root  /usr/local/bin/mindfulity 

# run on reboot
@reboot /usr/local/bin/mindfulity
